# x86/tags/CRUX-2_1/contrib/slim/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/slim
-rw-r--r--	root/root	etc/slim.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/slim
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/slim/
drwxr-xr-x	root/root	usr/share/slim/themes/
drwxr-xr-x	root/root	usr/share/slim/themes/curvy/
-rw-r--r--	root/root	usr/share/slim/themes/curvy/background.png
-rw-r--r--	root/root	usr/share/slim/themes/curvy/panel.png
-rw-r--r--	root/root	usr/share/slim/themes/curvy/slim.theme
drwxr-xr-x	root/root	usr/share/slim/themes/default/
-rw-r--r--	root/root	usr/share/slim/themes/default/background.png
-rw-r--r--	root/root	usr/share/slim/themes/default/panel.png
-rw-r--r--	root/root	usr/share/slim/themes/default/slim.theme
drwxr-xr-x	root/root	usr/share/slim/themes/leaves/
-rw-r--r--	root/root	usr/share/slim/themes/leaves/background.jpg
-rw-r--r--	root/root	usr/share/slim/themes/leaves/panel.png
-rw-r--r--	root/root	usr/share/slim/themes/leaves/slim.theme
drwxr-xr-x	root/root	usr/share/slim/themes/pirate/
-rw-r--r--	root/root	usr/share/slim/themes/pirate/background.png
-rw-r--r--	root/root	usr/share/slim/themes/pirate/panel.png
-rw-r--r--	root/root	usr/share/slim/themes/pirate/panel2.png
-rw-r--r--	root/root	usr/share/slim/themes/pirate/slim.theme
