# x86/tags/CRUX-2_1/contrib/hotplug-perl/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/hotplug-perl/
-rw-r--r--	root/root	etc/hotplug-perl/Hotplug.pm
drwxr-xr-x	root/root	etc/hotplug-perl/Hotplug/
-rw-r--r--	root/root	etc/hotplug-perl/Hotplug/IDE.pm
-rw-r--r--	root/root	etc/hotplug-perl/Hotplug/IEEE1394.pm
-rw-r--r--	root/root	etc/hotplug-perl/Hotplug/INPUT.pm
-rw-r--r--	root/root	etc/hotplug-perl/Hotplug/ISAPNP.pm
-rw-r--r--	root/root	etc/hotplug-perl/Hotplug/PCI.pm
-rw-r--r--	root/root	etc/hotplug-perl/Hotplug/SCSI.pm
-rw-r--r--	root/root	etc/hotplug-perl/Hotplug/USB.pm
-rw-r--r--	root/root	etc/hotplug-perl/blacklist
drwxr-xr-x	root/root	etc/hotplug-perl/blacklist.d/
-rwxr-xr-x	root/root	etc/hotplug-perl/list-detected-modules
-rwxr-xr-x	root/root	etc/hotplug-perl/start-hotplug
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/hotplug-perl
drwxr-xr-x	root/root	sbin/
-rwxr-xr-x	root/root	sbin/hotplug-perl
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man3/
-rw-r--r--	root/root	usr/share/man/man3/Hotplug.3.gz
-rw-r--r--	root/root	usr/share/man/man3/Hotplug::IDE.3.gz
-rw-r--r--	root/root	usr/share/man/man3/Hotplug::IEEE1394.3.gz
-rw-r--r--	root/root	usr/share/man/man3/Hotplug::INPUT.3.gz
-rw-r--r--	root/root	usr/share/man/man3/Hotplug::ISAPNP.3.gz
-rw-r--r--	root/root	usr/share/man/man3/Hotplug::PCI.3.gz
-rw-r--r--	root/root	usr/share/man/man3/Hotplug::SCSI.3.gz
-rw-r--r--	root/root	usr/share/man/man3/Hotplug::USB.3.gz
