# x86/tags/CRUX-2_1/contrib/dhcp/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/dhcp/
-rw-r-----	root/root	etc/dhcp/dhcpd.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/dhcpd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/omshell
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/dhcp/
-rw-r--r--	root/root	usr/include/dhcp/dhcpctl.h
drwxr-xr-x	root/root	usr/include/dhcp/isc-dhcp/
-rw-r--r--	root/root	usr/include/dhcp/isc-dhcp/boolean.h
-rw-r--r--	root/root	usr/include/dhcp/isc-dhcp/dst.h
-rw-r--r--	root/root	usr/include/dhcp/isc-dhcp/int.h
-rw-r--r--	root/root	usr/include/dhcp/isc-dhcp/lang.h
-rw-r--r--	root/root	usr/include/dhcp/isc-dhcp/list.h
-rw-r--r--	root/root	usr/include/dhcp/isc-dhcp/result.h
-rw-r--r--	root/root	usr/include/dhcp/isc-dhcp/types.h
drwxr-xr-x	root/root	usr/include/dhcp/omapip/
-rw-r--r--	root/root	usr/include/dhcp/omapip/alloc.h
-rw-r--r--	root/root	usr/include/dhcp/omapip/buffer.h
-rw-r--r--	root/root	usr/include/dhcp/omapip/omapip.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/dhcp/
-rw-r--r--	root/root	usr/lib/dhcp/libdhcpctl.a
-rw-r--r--	root/root	usr/lib/dhcp/libomapi.a
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rwxr-xr-x	root/root	usr/man/man1/omshell.1.gz
drwxr-xr-x	root/root	usr/man/man3/
-rwxr-xr-x	root/root	usr/man/man3/dhcpctl.3.gz
-rwxr-xr-x	root/root	usr/man/man3/omapi.3.gz
-rwxr-xr-x	root/root	usr/man/man3/omshell.3.gz
drwxr-xr-x	root/root	usr/man/man5/
-rwxr-xr-x	root/root	usr/man/man5/dhcp-eval.5.gz
-rwxr-xr-x	root/root	usr/man/man5/dhcp-options.5.gz
-rwxr-xr-x	root/root	usr/man/man5/dhcpd.conf.5.gz
-rwxr-xr-x	root/root	usr/man/man5/dhcpd.leases.5.gz
drwxr-xr-x	root/root	usr/man/man8/
-rwxr-xr-x	root/root	usr/man/man8/dhcpd.8.gz
-rwxr-xr-x	root/root	usr/man/man8/dhcrelay.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/dhcpd
-rwxr-xr-x	root/root	usr/sbin/dhcrelay
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/dhcp/
-rw-r--r--	root/root	var/lib/dhcp/dhcpd.leases (EMPTY)
drwxr-xr-x	root/root	var/run/
drwxr-xr-x	root/root	var/run/dhcp/
