# sparc/tags/CRUX-2_1/contrib/ntp/.footprint
drwxr-xr-x	peter/users	etc/
drwxr-xr-x	peter/users	etc/ntp/
-rw-r--r--	peter/users	etc/ntp/ntp.conf
drwxr-xr-x	peter/users	etc/rc.d/
-rwxr-xr-x	peter/users	etc/rc.d/ntpd
drwxr-xr-x	peter/users	usr/
drwxr-xr-x	peter/users	usr/bin/
-rwxr-xr-x	peter/users	usr/bin/ntp-keygen
-rwxr-xr-x	peter/users	usr/bin/ntp-wait
-rwxr-xr-x	peter/users	usr/bin/ntpd
-rwxr-xr-x	peter/users	usr/bin/ntpdate
-rwxr-xr-x	peter/users	usr/bin/ntpdc
-rwxr-xr-x	peter/users	usr/bin/ntpq
-rwxr-xr-x	peter/users	usr/bin/ntptime
-rwxr-xr-x	peter/users	usr/bin/ntptrace
-rwxr-xr-x	peter/users	usr/bin/sntp
-rwxr-xr-x	peter/users	usr/bin/tickadj
drwxr-xr-x	peter/users	usr/man/
drwxr-xr-x	peter/users	usr/man/man5/
-rw-r--r--	peter/users	usr/man/man5/ntp.conf.5.gz
-rw-r--r--	peter/users	usr/man/man5/ntp.keys.5.gz
drwxr-xr-x	peter/users	usr/man/man8/
-rw-r--r--	peter/users	usr/man/man8/ntp-genkeys.8.gz
-rw-r--r--	peter/users	usr/man/man8/ntpd.8.gz
-rw-r--r--	peter/users	usr/man/man8/ntpdate.8.gz
-rw-r--r--	peter/users	usr/man/man8/ntpdc.8.gz
-rw-r--r--	peter/users	usr/man/man8/ntpq.8.gz
-rw-r--r--	peter/users	usr/man/man8/ntptime.8.gz
-rw-r--r--	peter/users	usr/man/man8/ntptrace.8.gz
drwxr-xr-x	peter/users	var/
drwxr-xr-x	peter/users	var/lib/
drwxr-xr-x	peter/users	var/lib/ntp/
-rw-r--r--	peter/users	var/lib/ntp/ntp.drift
drwxr-xr-x	peter/users	var/lib/ntp/stats/
drwxr-xr-x	peter/users	var/log/
-rw-r--r--	peter/users	var/log/ntp.log (EMPTY)
drwxr-xr-x	peter/users	var/run/
drwxr-xr-x	peter/users	var/run/ntp/
