# x86/tags/CRUX-2_1/contrib/jakarta-tomcat/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/tomcat
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/tomcat/
-rw-r--r--	root/root	usr/lib/tomcat/NOTICE
drwxr-xr-x	root/root	usr/lib/tomcat/bin/
-rw-r--r--	root/root	usr/lib/tomcat/bin/LauncherBootstrap.class
-rw-r--r--	root/root	usr/lib/tomcat/bin/bootstrap.jar
-rw-r--r--	root/root	usr/lib/tomcat/bin/catalina.50.sh
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/catalina.sh
-rw-r--r--	root/root	usr/lib/tomcat/bin/catalina.xml
-rw-r--r--	root/root	usr/lib/tomcat/bin/commons-daemon.jar
-rw-r--r--	root/root	usr/lib/tomcat/bin/commons-launcher.jar
-rw-r--r--	root/root	usr/lib/tomcat/bin/commons-logging-api.jar
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/digest.sh
-rw-r--r--	root/root	usr/lib/tomcat/bin/jmx.jar
-rw-r--r--	root/root	usr/lib/tomcat/bin/jsvc.tar.gz
-rw-r--r--	root/root	usr/lib/tomcat/bin/launcher.properties
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/setclasspath.sh
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/shutdown-using-launcher.sh
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/shutdown.sh
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/startup-using-launcher.sh
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/startup.sh
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/tool-wrapper-using-launcher.sh
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/tool-wrapper.sh
-rw-r--r--	root/root	usr/lib/tomcat/bin/version.sh
drwxr-xr-x	root/root	usr/lib/tomcat/common/
drwxr-xr-x	root/root	usr/lib/tomcat/common/classes/
drwxr-xr-x	root/root	usr/lib/tomcat/common/endorsed/
-rw-r--r--	root/root	usr/lib/tomcat/common/endorsed/xercesImpl.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/endorsed/xml-apis.jar
drwxr-xr-x	root/root	usr/lib/tomcat/common/lib/
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/ant-launcher.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/ant.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/commons-collections-3.1.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/commons-dbcp-1.2.1.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/commons-el.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/commons-pool-1.2.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/jasper-compiler.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/jasper-runtime.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/jsp-api.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/naming-common.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/naming-factory.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/naming-java.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/naming-resources.jar
-rw-r--r--	root/root	usr/lib/tomcat/common/lib/servlet-api.jar
drwxr-xr-x	root/root	usr/lib/tomcat/conf/
drwxr-xr-x	root/root	usr/lib/tomcat/conf/Catalina/
drwxr-xr-x	root/root	usr/lib/tomcat/conf/Catalina/localhost/
-rw-------	root/root	usr/lib/tomcat/conf/Catalina/localhost/admin.xml
-rw-------	root/root	usr/lib/tomcat/conf/Catalina/localhost/balancer.xml
-rw-------	root/root	usr/lib/tomcat/conf/Catalina/localhost/manager.xml
-rw-------	root/root	usr/lib/tomcat/conf/catalina.policy
-rw-------	root/root	usr/lib/tomcat/conf/catalina.properties
-rw-------	root/root	usr/lib/tomcat/conf/jk2.properties
-rw-------	root/root	usr/lib/tomcat/conf/server-minimal.xml
-rw-------	root/root	usr/lib/tomcat/conf/server.xml
-rw-------	root/root	usr/lib/tomcat/conf/tomcat-users.xml
-rw-------	root/root	usr/lib/tomcat/conf/web.xml
drwxr-xr-x	root/root	usr/lib/tomcat/logs/
drwxr-xr-x	root/root	usr/lib/tomcat/server/
drwxr-xr-x	root/root	usr/lib/tomcat/server/classes/
drwxr-xr-x	root/root	usr/lib/tomcat/server/lib/
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/catalina-ant.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/catalina-cluster.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/catalina-i18n-es.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/catalina-i18n-fr.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/catalina-i18n-ja.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/catalina-optional.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/catalina.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/commons-beanutils.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/commons-digester.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/commons-fileupload-1.0.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/commons-modeler.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/jakarta-regexp-1.3.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/jkconfig.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/jkshm.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/servlets-cgi.renametojar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/servlets-common.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/servlets-default.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/servlets-invoker.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/servlets-ssi.renametojar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/servlets-webdav.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/tomcat-coyote.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/tomcat-http11.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/tomcat-jk.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/tomcat-jk2.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/tomcat-jni.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/lib/tomcat-util.jar
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/WEB-INF/
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/WEB-INF/controls.tld
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/WEB-INF/lib/
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/WEB-INF/lib/catalina-admin.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/WEB-INF/lib/struts.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/WEB-INF/struts-bean.tld
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/WEB-INF/struts-config.xml
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/WEB-INF/struts-html.tld
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/WEB-INF/struts-logic.tld
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/WEB-INF/web.xml
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/admin.css
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/admin.xml
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/connector/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/context/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/defaultcontext/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/host/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/images/
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/BlueTile.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Connector.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Context.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Datasource.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/DefaultContext.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/EnvironmentEntries.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Groups.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Host.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Logger.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Login.jpg
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/LoginBackgroundTile.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Mailsession.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/PaperTexture.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Realm.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/ResourceLink.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Roles.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Server.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Service.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Thumbs.db
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/TomcatBanner.jpg
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Users.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/Valve.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/folder_16_pad.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/handledownlast.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/handledownmiddle.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/handlerightlast.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/handlerightmiddle.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/linelastnode.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/linemiddlenode.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/images/linevertical.gif
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/logger/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/realm/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/resources/
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/resources/dataSources.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/resources/envEntries.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/resources/listDataSources.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/resources/listEnvEntries.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/resources/listMailSessions.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/resources/listResourceLinks.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/resources/listUserDatabases.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/resources/mailSessions.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/resources/resourceLinks.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/resources/userDatabases.jspf
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/server/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/service/
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/tree-control-test.css
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/users/
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/users/groups.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/users/listGroups.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/users/listRoles.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/users/listUsers.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/users/roles.jspf
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/admin/users/users.jspf
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/admin/valve/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/manager/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/manager/WEB-INF/
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/manager/WEB-INF/lib/
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/WEB-INF/lib/catalina-manager.jar
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/WEB-INF/web.xml
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/html-manager-howto.html
drwxr-xr-x	root/root	usr/lib/tomcat/server/webapps/manager/images/
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/add.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/code.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/design.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/docs.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/fix.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/jakarta-logo.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/tomcat.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/update.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/images/void.gif
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/manager-howto.html
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/manager.xml
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/status.xsd
-rw-r--r--	root/root	usr/lib/tomcat/server/webapps/manager/xform.xsl
drwxr-xr-x	root/root	usr/lib/tomcat/shared/
drwxr-xr-x	root/root	usr/lib/tomcat/shared/classes/
drwxr-xr-x	root/root	usr/lib/tomcat/shared/lib/
drwxr-xr-x	root/root	usr/lib/tomcat/temp/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/ROOT/
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/RELEASE-NOTES.txt
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/ROOT/WEB-INF/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/ROOT/WEB-INF/classes/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/ROOT/WEB-INF/classes/org/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/ROOT/WEB-INF/classes/org/apache/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/ROOT/WEB-INF/classes/org/apache/jsp/
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/WEB-INF/classes/org/apache/jsp/index_jsp.class
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/ROOT/WEB-INF/lib/
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/WEB-INF/lib/catalina-root.jar
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/WEB-INF/web.xml
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/index.jsp
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/jakarta-banner.gif
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/tomcat-power.gif
-rw-r--r--	root/root	usr/lib/tomcat/webapps/ROOT/tomcat.gif
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/balancer/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/balancer/WEB-INF/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/balancer/WEB-INF/config/
-rw-r--r--	root/root	usr/lib/tomcat/webapps/balancer/WEB-INF/config/rules.xml
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/balancer/WEB-INF/lib/
-rw-r--r--	root/root	usr/lib/tomcat/webapps/balancer/WEB-INF/lib/catalina-balancer.jar
-rw-r--r--	root/root	usr/lib/tomcat/webapps/balancer/WEB-INF/web.xml
-rw-r--r--	root/root	usr/lib/tomcat/webapps/balancer/balancer.xml
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/balancer/images/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/webdav/
drwxr-xr-x	root/root	usr/lib/tomcat/webapps/webdav/WEB-INF/
-rw-r--r--	root/root	usr/lib/tomcat/webapps/webdav/WEB-INF/web.xml
-rw-r--r--	root/root	usr/lib/tomcat/webapps/webdav/index.html
-rw-r--r--	root/root	usr/lib/tomcat/webapps/webdav/tomcat-power.gif
-rw-r--r--	root/root	usr/lib/tomcat/webapps/webdav/tomcat.gif
drwxr-xr-x	root/root	usr/lib/tomcat/work/
