# x86/tags/CRUX-2_1/contrib/ipsec-tools/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/racoon/
drwxr-xr-x	root/root	etc/racoon/certs/
-rw-r--r--	root/root	etc/racoon/certs/Makefile
drwxr-xr-x	root/root	etc/racoon/samples/
-rw-r--r--	root/root	etc/racoon/samples/psk.txt.sample
-rw-r--r--	root/root	etc/racoon/samples/racoon.conf.sample
-rw-r--r--	root/root	etc/racoon/samples/racoon.conf.sample-gssapi
-rw-r--r--	root/root	etc/racoon/samples/racoon.conf.sample-inherit
-rw-r--r--	root/root	etc/racoon/samples/racoon.conf.sample-natt
-rw-r--r--	root/root	etc/racoon/samples/racoon.conf.sample-plainrsa
-rw-------	root/root	etc/racoon/setkey.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/racoon
-rwxr-xr-x	root/root	etc/rc.d/setkey
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/libipsec/
-rw-r--r--	root/root	usr/include/libipsec/libpfkey.h
drwxr-xr-x	root/root	usr/lib/
-rw-r--r--	root/root	usr/lib/libipsec.a
-rwxr-xr-x	root/root	usr/lib/libipsec.la
lrwxrwxrwx	root/root	usr/lib/libipsec.so -> libipsec.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libipsec.so.0 -> libipsec.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libipsec.so.0.0.0
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man3/
-rw-r--r--	root/root	usr/man/man3/ipsec_set_policy.3.gz
-rw-r--r--	root/root	usr/man/man3/ipsec_strerror.3.gz
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/racoon.conf.5.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/plainrsa-gen.8.gz
-rw-r--r--	root/root	usr/man/man8/racoon.8.gz
-rw-r--r--	root/root	usr/man/man8/racoonctl.8.gz
-rw-r--r--	root/root	usr/man/man8/setkey.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/plainrsa-gen
-rwxr-xr-x	root/root	usr/sbin/racoon
-rwxr-xr-x	root/root	usr/sbin/racoonctl
-rwxr-xr-x	root/root	usr/sbin/setkey
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/racoon/
