# x86/tags/CRUX-2_1/contrib/shorewall/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-r-xr--r--	root/root	etc/rc.d/shorewall
drwxr-xr-x	root/root	etc/shorewall/
-rw-------	root/root	etc/shorewall/accounting
-rw-------	root/root	etc/shorewall/actions
-rw-------	root/root	etc/shorewall/blacklist
-rw-------	root/root	etc/shorewall/ecn
-rw-------	root/root	etc/shorewall/hosts
-rw-------	root/root	etc/shorewall/init
-rw-------	root/root	etc/shorewall/initdone
-rw-------	root/root	etc/shorewall/interfaces
-rw-------	root/root	etc/shorewall/ipsec
-rw-------	root/root	etc/shorewall/maclist
-rw-------	root/root	etc/shorewall/masq
-rw-------	root/root	etc/shorewall/modules
-rw-------	root/root	etc/shorewall/nat
-rw-------	root/root	etc/shorewall/netmap
-rw-------	root/root	etc/shorewall/params
-rw-------	root/root	etc/shorewall/policy
-rw-------	root/root	etc/shorewall/proxyarp
-rw-------	root/root	etc/shorewall/routestopped
-rw-------	root/root	etc/shorewall/rules
-rwxr--r--	root/root	etc/shorewall/shorewall.conf
-rw-------	root/root	etc/shorewall/start
-rw-------	root/root	etc/shorewall/stop
-rw-------	root/root	etc/shorewall/stopped
-rw-------	root/root	etc/shorewall/tcrules
-rw-------	root/root	etc/shorewall/tos
-rw-------	root/root	etc/shorewall/tunnels
-rwxr--r--	root/root	etc/shorewall/zones
drwxr-xr-x	root/root	sbin/
-r-xr--r--	root/root	sbin/shorewall
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/shorewall/
-rw-------	root/root	usr/share/shorewall/action.AllowAuth
-rw-------	root/root	usr/share/shorewall/action.AllowDNS
-rw-------	root/root	usr/share/shorewall/action.AllowFTP
-rw-------	root/root	usr/share/shorewall/action.AllowICMPs
-rw-------	root/root	usr/share/shorewall/action.AllowIMAP
-rw-------	root/root	usr/share/shorewall/action.AllowNNTP
-rw-------	root/root	usr/share/shorewall/action.AllowNTP
-rw-------	root/root	usr/share/shorewall/action.AllowPCA
-rw-------	root/root	usr/share/shorewall/action.AllowPOP3
-rw-------	root/root	usr/share/shorewall/action.AllowPing
-rw-------	root/root	usr/share/shorewall/action.AllowRdate
-rw-------	root/root	usr/share/shorewall/action.AllowSMB
-rw-------	root/root	usr/share/shorewall/action.AllowSMTP
-rw-------	root/root	usr/share/shorewall/action.AllowSNMP
-rw-------	root/root	usr/share/shorewall/action.AllowSSH
-rw-------	root/root	usr/share/shorewall/action.AllowTelnet
-rw-------	root/root	usr/share/shorewall/action.AllowTrcrt
-rw-------	root/root	usr/share/shorewall/action.AllowVNC
-rw-------	root/root	usr/share/shorewall/action.AllowVNCL
-rw-------	root/root	usr/share/shorewall/action.AllowWeb
-rw-------	root/root	usr/share/shorewall/action.Drop
-rw-------	root/root	usr/share/shorewall/action.DropDNSrep
-rw-------	root/root	usr/share/shorewall/action.DropPing
-rw-------	root/root	usr/share/shorewall/action.DropSMB
-rw-------	root/root	usr/share/shorewall/action.DropUPnP
-rw-------	root/root	usr/share/shorewall/action.Reject
-rw-------	root/root	usr/share/shorewall/action.RejectAuth
-rw-------	root/root	usr/share/shorewall/action.RejectSMB
-rw-------	root/root	usr/share/shorewall/action.template
-rw-------	root/root	usr/share/shorewall/actions.std
-rw-------	root/root	usr/share/shorewall/bogons
-rw-------	root/root	usr/share/shorewall/configpath
-r-xr--r--	root/root	usr/share/shorewall/firewall
-r--r--r--	root/root	usr/share/shorewall/functions
-r-xr--r--	root/root	usr/share/shorewall/help
-rw-------	root/root	usr/share/shorewall/rfc1918
-rw-r--r--	root/root	usr/share/shorewall/version
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/shorewall/
