# x86/tags/CRUX-2_1/contrib/bind/.footprint
drwxr-xr-x	root/root	etc/
-rw-r-----	root/named	etc/named.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/named
-rw-r-----	root/named	etc/rndc.conf (EMPTY)
-rw-r-----	root/named	etc/rndc.key (EMPTY)
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/isc-config.sh
-rwxr-xr-x	root/root	usr/bin/nslookup
-rwxr-xr-x	root/root	usr/bin/nsupdate
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/bind9/
-rw-r--r--	root/root	usr/include/bind9/check.h
-rw-r--r--	root/root	usr/include/bind9/getaddresses.h
-rw-r--r--	root/root	usr/include/bind9/version.h
drwxr-xr-x	root/root	usr/include/dns/
-rw-r--r--	root/root	usr/include/dns/acl.h
-rw-r--r--	root/root	usr/include/dns/adb.h
-rw-r--r--	root/root	usr/include/dns/byaddr.h
-rw-r--r--	root/root	usr/include/dns/cache.h
-rw-r--r--	root/root	usr/include/dns/callbacks.h
-rw-r--r--	root/root	usr/include/dns/cert.h
-rw-r--r--	root/root	usr/include/dns/compress.h
-rw-r--r--	root/root	usr/include/dns/db.h
-rw-r--r--	root/root	usr/include/dns/dbiterator.h
-rw-r--r--	root/root	usr/include/dns/dbtable.h
-rw-r--r--	root/root	usr/include/dns/diff.h
-rw-r--r--	root/root	usr/include/dns/dispatch.h
-rw-r--r--	root/root	usr/include/dns/dnssec.h
-rw-r--r--	root/root	usr/include/dns/ds.h
-rw-r--r--	root/root	usr/include/dns/enumclass.h
-rw-r--r--	root/root	usr/include/dns/enumtype.h
-rw-r--r--	root/root	usr/include/dns/events.h
-rw-r--r--	root/root	usr/include/dns/fixedname.h
-rw-r--r--	root/root	usr/include/dns/journal.h
-rw-r--r--	root/root	usr/include/dns/keyflags.h
-rw-r--r--	root/root	usr/include/dns/keytable.h
-rw-r--r--	root/root	usr/include/dns/keyvalues.h
-rw-r--r--	root/root	usr/include/dns/lib.h
-rw-r--r--	root/root	usr/include/dns/log.h
-rw-r--r--	root/root	usr/include/dns/master.h
-rw-r--r--	root/root	usr/include/dns/masterdump.h
-rw-r--r--	root/root	usr/include/dns/message.h
-rw-r--r--	root/root	usr/include/dns/name.h
-rw-r--r--	root/root	usr/include/dns/ncache.h
-rw-r--r--	root/root	usr/include/dns/nsec.h
-rw-r--r--	root/root	usr/include/dns/peer.h
-rw-r--r--	root/root	usr/include/dns/portlist.h
-rw-r--r--	root/root	usr/include/dns/rbt.h
-rw-r--r--	root/root	usr/include/dns/rcode.h
-rw-r--r--	root/root	usr/include/dns/rdata.h
-rw-r--r--	root/root	usr/include/dns/rdataclass.h
-rw-r--r--	root/root	usr/include/dns/rdatalist.h
-rw-r--r--	root/root	usr/include/dns/rdataset.h
-rw-r--r--	root/root	usr/include/dns/rdatasetiter.h
-rw-r--r--	root/root	usr/include/dns/rdataslab.h
-rw-r--r--	root/root	usr/include/dns/rdatastruct.h
-rw-r--r--	root/root	usr/include/dns/rdatatype.h
-rw-r--r--	root/root	usr/include/dns/request.h
-rw-r--r--	root/root	usr/include/dns/resolver.h
-rw-r--r--	root/root	usr/include/dns/result.h
-rw-r--r--	root/root	usr/include/dns/rootns.h
-rw-r--r--	root/root	usr/include/dns/sdb.h
-rw-r--r--	root/root	usr/include/dns/secalg.h
-rw-r--r--	root/root	usr/include/dns/secproto.h
-rw-r--r--	root/root	usr/include/dns/soa.h
-rw-r--r--	root/root	usr/include/dns/ssu.h
-rw-r--r--	root/root	usr/include/dns/tcpmsg.h
-rw-r--r--	root/root	usr/include/dns/time.h
-rw-r--r--	root/root	usr/include/dns/tkey.h
-rw-r--r--	root/root	usr/include/dns/tsig.h
-rw-r--r--	root/root	usr/include/dns/ttl.h
-rw-r--r--	root/root	usr/include/dns/types.h
-rw-r--r--	root/root	usr/include/dns/validator.h
-rw-r--r--	root/root	usr/include/dns/version.h
-rw-r--r--	root/root	usr/include/dns/view.h
-rw-r--r--	root/root	usr/include/dns/xfrin.h
-rw-r--r--	root/root	usr/include/dns/zone.h
-rw-r--r--	root/root	usr/include/dns/zonekey.h
-rw-r--r--	root/root	usr/include/dns/zt.h
drwxr-xr-x	root/root	usr/include/dst/
-rw-r--r--	root/root	usr/include/dst/dst.h
-rw-r--r--	root/root	usr/include/dst/lib.h
-rw-r--r--	root/root	usr/include/dst/result.h
drwxr-xr-x	root/root	usr/include/isc/
-rw-r--r--	root/root	usr/include/isc/app.h
-rw-r--r--	root/root	usr/include/isc/assertions.h
-rw-r--r--	root/root	usr/include/isc/base64.h
-rw-r--r--	root/root	usr/include/isc/bitstring.h
-rw-r--r--	root/root	usr/include/isc/boolean.h
-rw-r--r--	root/root	usr/include/isc/buffer.h
-rw-r--r--	root/root	usr/include/isc/bufferlist.h
-rw-r--r--	root/root	usr/include/isc/commandline.h
-rw-r--r--	root/root	usr/include/isc/condition.h
-rw-r--r--	root/root	usr/include/isc/dir.h
-rw-r--r--	root/root	usr/include/isc/entropy.h
-rw-r--r--	root/root	usr/include/isc/error.h
-rw-r--r--	root/root	usr/include/isc/event.h
-rw-r--r--	root/root	usr/include/isc/eventclass.h
-rw-r--r--	root/root	usr/include/isc/file.h
-rw-r--r--	root/root	usr/include/isc/formatcheck.h
-rw-r--r--	root/root	usr/include/isc/fsaccess.h
-rw-r--r--	root/root	usr/include/isc/heap.h
-rw-r--r--	root/root	usr/include/isc/hex.h
-rw-r--r--	root/root	usr/include/isc/hmacmd5.h
-rw-r--r--	root/root	usr/include/isc/int.h
-rw-r--r--	root/root	usr/include/isc/interfaceiter.h
-rw-r--r--	root/root	usr/include/isc/lang.h
-rw-r--r--	root/root	usr/include/isc/lex.h
-rw-r--r--	root/root	usr/include/isc/lfsr.h
-rw-r--r--	root/root	usr/include/isc/lib.h
-rw-r--r--	root/root	usr/include/isc/list.h
-rw-r--r--	root/root	usr/include/isc/log.h
-rw-r--r--	root/root	usr/include/isc/magic.h
-rw-r--r--	root/root	usr/include/isc/md5.h
-rw-r--r--	root/root	usr/include/isc/mem.h
-rw-r--r--	root/root	usr/include/isc/msgcat.h
-rw-r--r--	root/root	usr/include/isc/msgs.h
-rw-r--r--	root/root	usr/include/isc/mutex.h
-rw-r--r--	root/root	usr/include/isc/mutexblock.h
-rw-r--r--	root/root	usr/include/isc/net.h
-rw-r--r--	root/root	usr/include/isc/netaddr.h
-rw-r--r--	root/root	usr/include/isc/netdb.h
-rw-r--r--	root/root	usr/include/isc/offset.h
-rw-r--r--	root/root	usr/include/isc/once.h
-rw-r--r--	root/root	usr/include/isc/ondestroy.h
-rw-r--r--	root/root	usr/include/isc/os.h
-rw-r--r--	root/root	usr/include/isc/parseint.h
-rw-r--r--	root/root	usr/include/isc/platform.h
-rw-r--r--	root/root	usr/include/isc/print.h
-rw-r--r--	root/root	usr/include/isc/quota.h
-rw-r--r--	root/root	usr/include/isc/random.h
-rw-r--r--	root/root	usr/include/isc/ratelimiter.h
-rw-r--r--	root/root	usr/include/isc/refcount.h
-rw-r--r--	root/root	usr/include/isc/region.h
-rw-r--r--	root/root	usr/include/isc/resource.h
-rw-r--r--	root/root	usr/include/isc/result.h
-rw-r--r--	root/root	usr/include/isc/resultclass.h
-rw-r--r--	root/root	usr/include/isc/rwlock.h
-rw-r--r--	root/root	usr/include/isc/serial.h
-rw-r--r--	root/root	usr/include/isc/sha1.h
-rw-r--r--	root/root	usr/include/isc/sockaddr.h
-rw-r--r--	root/root	usr/include/isc/socket.h
-rw-r--r--	root/root	usr/include/isc/stdio.h
-rw-r--r--	root/root	usr/include/isc/stdlib.h
-rw-r--r--	root/root	usr/include/isc/stdtime.h
-rw-r--r--	root/root	usr/include/isc/string.h
-rw-r--r--	root/root	usr/include/isc/symtab.h
-rw-r--r--	root/root	usr/include/isc/syslog.h
-rw-r--r--	root/root	usr/include/isc/task.h
-rw-r--r--	root/root	usr/include/isc/taskpool.h
-rw-r--r--	root/root	usr/include/isc/thread.h
-rw-r--r--	root/root	usr/include/isc/time.h
-rw-r--r--	root/root	usr/include/isc/timer.h
-rw-r--r--	root/root	usr/include/isc/types.h
-rw-r--r--	root/root	usr/include/isc/util.h
-rw-r--r--	root/root	usr/include/isc/version.h
drwxr-xr-x	root/root	usr/include/isccc/
-rw-r--r--	root/root	usr/include/isccc/alist.h
-rw-r--r--	root/root	usr/include/isccc/base64.h
-rw-r--r--	root/root	usr/include/isccc/cc.h
-rw-r--r--	root/root	usr/include/isccc/ccmsg.h
-rw-r--r--	root/root	usr/include/isccc/events.h
-rw-r--r--	root/root	usr/include/isccc/lib.h
-rw-r--r--	root/root	usr/include/isccc/result.h
-rw-r--r--	root/root	usr/include/isccc/sexpr.h
-rw-r--r--	root/root	usr/include/isccc/symtab.h
-rw-r--r--	root/root	usr/include/isccc/symtype.h
-rw-r--r--	root/root	usr/include/isccc/types.h
-rw-r--r--	root/root	usr/include/isccc/util.h
-rw-r--r--	root/root	usr/include/isccc/version.h
drwxr-xr-x	root/root	usr/include/isccfg/
-rw-r--r--	root/root	usr/include/isccfg/cfg.h
-rw-r--r--	root/root	usr/include/isccfg/grammar.h
-rw-r--r--	root/root	usr/include/isccfg/log.h
-rw-r--r--	root/root	usr/include/isccfg/namedconf.h
-rw-r--r--	root/root	usr/include/isccfg/version.h
drwxr-xr-x	root/root	usr/include/lwres/
-rw-r--r--	root/root	usr/include/lwres/context.h
-rw-r--r--	root/root	usr/include/lwres/int.h
-rw-r--r--	root/root	usr/include/lwres/ipv6.h
-rw-r--r--	root/root	usr/include/lwres/lang.h
-rw-r--r--	root/root	usr/include/lwres/list.h
-rw-r--r--	root/root	usr/include/lwres/lwbuffer.h
-rw-r--r--	root/root	usr/include/lwres/lwpacket.h
-rw-r--r--	root/root	usr/include/lwres/lwres.h
-rw-r--r--	root/root	usr/include/lwres/net.h
-rw-r--r--	root/root	usr/include/lwres/netdb.h
-rw-r--r--	root/root	usr/include/lwres/platform.h
-rw-r--r--	root/root	usr/include/lwres/result.h
-rw-r--r--	root/root	usr/include/lwres/version.h
drwxr-xr-x	root/root	usr/lib/
-rw-r--r--	root/root	usr/lib/libbind9.a
-rw-r--r--	root/root	usr/lib/libdns.a
-rw-r--r--	root/root	usr/lib/libisc.a
-rw-r--r--	root/root	usr/lib/libisccc.a
-rw-r--r--	root/root	usr/lib/libisccfg.a
-rw-r--r--	root/root	usr/lib/liblwres.a
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man3/
-rw-r--r--	root/root	usr/man/man3/lwres.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_addr_parse.3.gz -> lwres_resutil.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_add.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_back.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_clear.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_first.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_forward.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_getmem.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_getuint16.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_getuint32.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_getuint8.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_init.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_invalidate.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_putmem.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_putuint16.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_putuint32.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_putuint8.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_buffer_subtract.3.gz -> lwres_buffer.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_conf_clear.3.gz -> lwres_config.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_conf_get.3.gz -> lwres_config.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_conf_init.3.gz -> lwres_config.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_conf_parse.3.gz -> lwres_config.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_conf_print.3.gz -> lwres_config.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_config.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_context.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_context_allocmem.3.gz -> lwres_context.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_context_create.3.gz -> lwres_context.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_context_destroy.3.gz -> lwres_context.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_context_freemem.3.gz -> lwres_context.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_context_initserial.3.gz -> lwres_context.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_context_nextserial.3.gz -> lwres_context.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_context_sendrecv.3.gz -> lwres_context.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_endhostent.3.gz -> lwres_gethostent.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_endhostent_r.3.gz -> lwres_gethostent.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_freeaddrinfo.3.gz -> lwres_getaddrinfo.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_freehostent.3.gz -> lwres_getipnode.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_gabn.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gabnrequest_free.3.gz -> lwres_gabn.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gabnrequest_parse.3.gz -> lwres_gabn.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gabnrequest_render.3.gz -> lwres_gabn.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gabnresponse_free.3.gz -> lwres_gabn.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gabnresponse_parse.3.gz -> lwres_gabn.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gabnresponse_render.3.gz -> lwres_gabn.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_gai_strerror.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_getaddrinfo.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_getaddrsbyname.3.gz -> lwres_resutil.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gethostbyaddr.3.gz -> lwres_gethostent.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gethostbyaddr_r.3.gz -> lwres_gethostent.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gethostbyname.3.gz -> lwres_gethostent.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gethostbyname2.3.gz -> lwres_gethostent.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gethostbyname_r.3.gz -> lwres_gethostent.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_gethostent.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gethostent_r.3.gz -> lwres_gethostent.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_getipnode.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_getipnodebyaddr.3.gz -> lwres_getipnode.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_getipnodebyname.3.gz -> lwres_getipnode.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_getnamebyaddr.3.gz -> lwres_resutil.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_getnameinfo.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_getrrsetbyname.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_gnba.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gnbarequest_free.3.gz -> lwres_gnba.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gnbarequest_parse.3.gz -> lwres_gnba.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gnbarequest_render.3.gz -> lwres_gnba.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gnbaresponse_free.3.gz -> lwres_gnba.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gnbaresponse_parse.3.gz -> lwres_gnba.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_gnbaresponse_render.3.gz -> lwres_gnba.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_herror.3.gz -> lwres_hstrerror.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_hstrerror.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_inetntop.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_lwpacket_parseheader.3.gz -> lwres_packet.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_lwpacket_renderheader.3.gz -> lwres_packet.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_net_ntop.3.gz -> lwres_inetntop.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_noop.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_nooprequest_free.3.gz -> lwres_noop.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_nooprequest_parse.3.gz -> lwres_noop.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_nooprequest_render.3.gz -> lwres_noop.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_noopresponse_free.3.gz -> lwres_noop.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_noopresponse_parse.3.gz -> lwres_noop.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_noopresponse_render.3.gz -> lwres_noop.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_packet.3.gz
-rw-r--r--	root/root	usr/man/man3/lwres_resutil.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_sethostent.3.gz -> lwres_gethostent.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_sethostent_r.3.gz -> lwres_gethostent.3.gz
lrwxrwxrwx	root/root	usr/man/man3/lwres_string_parse.3.gz -> lwres_resutil.3.gz
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/named.conf.5.gz
-rw-r--r--	root/root	usr/man/man5/rndc.conf.5.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/dnssec-keygen.8.gz
-rw-r--r--	root/root	usr/man/man8/dnssec-signzone.8.gz
-rw-r--r--	root/root	usr/man/man8/lwresd.8.gz
-rw-r--r--	root/root	usr/man/man8/named-checkconf.8.gz
-rw-r--r--	root/root	usr/man/man8/named-checkzone.8.gz
-rw-r--r--	root/root	usr/man/man8/named.8.gz
-rw-r--r--	root/root	usr/man/man8/nsupdate.8.gz
-rw-r--r--	root/root	usr/man/man8/rndc-confgen.8.gz
-rw-r--r--	root/root	usr/man/man8/rndc.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/dnssec-keygen
-rwxr-xr-x	root/root	usr/sbin/dnssec-signzone
-rwxr-xr-x	root/root	usr/sbin/lwresd
-rwxr-xr-x	root/root	usr/sbin/named
-rwxr-xr-x	root/root	usr/sbin/named-checkconf
-rwxr-xr-x	root/root	usr/sbin/named-checkzone
-rwxr-xr-x	root/root	usr/sbin/rndc
-rwxr-xr-x	root/root	usr/sbin/rndc-confgen
drwxr-xr-x	root/root	var/
drwxrwx--T	root/named	var/named/
-rw-r-----	root/named	var/named/db.0.0.127
-rw-r-----	root/named	var/named/db.ip6.arpa
-rw-r-----	root/named	var/named/db.localhost
-rw-r-----	root/named	var/named/named.root
drwxr-xr-x	root/root	var/run/
drwxrwx---	root/named	var/run/named/
-rw-r--r--	named/named	var/run/named/named.pid (EMPTY)
